# automation/systemd/gestion-telegram-bot.service
[Unit]
Description=Automation Telegram Bot
After=network.target gestion-scheduler.service
Wants=network.target

[Service]
Type=simple
User=andres-osorio
Group=andres-osorio
WorkingDirectory=/home/andres-osorio/server_dato/automation
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/andres-osorio/server_dato/automation/.env
ExecStart=/home/andres-osorio/server_dato/automation/venv/bin/python /home/andres-osorio/server_dato/automation/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
